# files/vpnd.in
#!/bin/sh
#

# PROVIDE: vpnd
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name=vpnd
rcvar=vpnd_enable

load_rc_config $name

: ${vpnd_enable="NO"}
: ${vpnd_pidfile="/var/run/${name}.pid"}
: ${vpnd_log="/var/log/vpnd.log"}

command="/usr/local/sbin/${name}"
start_cmd="/usr/sbin/daemon -p ${vpnd_pidfile} -f -o ${vpnd_log} ${vpnd_command}"
pidfile="${vpnd_pidfile}"

run_rc_command $*
